# systemd/coclaude-xiaog.service
[Unit]
Description=蟆秀 coclaude 蝗｢髦溷刻菴懈恪蜉｡
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash /home/nvidia/localmodel/coclaude-service.sh
ExecStop=/usr/bin/tmux kill-session -t coclaude-team
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
